# npm/docker-autoclean.service
[Unit]
Description=Automatic Docker cleanup

[Service]
Type=oneshot
ExecStart=/usr/local/bin/docker-autoclean.sh